# server/chess-mcp.service
[Unit]
Description=Chess MCP Analysis Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/chess-mcp/server
ExecStart=/usr/bin/env uv run chess_mcp.py
Restart=on-failure
Environment=FASTMCP_HOST=0.0.0.0
Environment=FASTMCP_PORT=8000
Environment=STOCKFISH_PATH=/usr/bin/stockfish
Environment=ANALYSIS_DEPTH=18

[Install]
WantedBy=default.target
